# backend/services/automation-engine/supervisor.conf
[program:automation-engine]
process_name=%(program_name)s
command=python -u main.py
directory=/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/automation-engine.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=10
priority=999

