# backend/services/automation-engine/supervisor.conf
[program:automation-engine]
process_name=%(program_name)s
command=python -u main.py
directory=/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=1
redirect_stderr=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopwaitsecs=10
priority=999
